# cron.d/closelunch
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
0 13 * * * root /usr/local/bin/closefood.sh > /tmp/stdoutf 2>/tmp/errorf
