# cron.d/closeduckdoor
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
30 18 * * * root /usr/local/bin/closedoor.sh > /tmp/stdout 2>/tmp/error
