# cron.d/openduckdoor
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
30 7 * * * root /usr/local/bin/opendoor.sh > /tmp/stdout 2>/tmp/error
