# cron.d/openlunch
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
0 11 * * * root /usr/local/bin/openfood.sh > /tmp/stdoutf 2>/tmp/errorf
