# cron.d/opendinner
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
29 19 * * * root /usr/local/bin/openfood.sh > /tmp/stdoutf 2>/tmp/errorf
